# systemd/usr-share-fex\x2demu\x2drootfs\x2dlayers-gentoo-layers-00\x2dbase.mount
[Unit]
Description=Fex rootfs base layer

[Mount]
What=/usr/share/fex-emu-rootfs-layers/gentoo/images/00-base.sqfs
Where=/usr/share/fex-emu-rootfs-layers/gentoo/layers/00-base/
