# systemd/usr-share-fex\x2demu\x2drootfs\x2dlayers-gentoo-layers-30\x2dmesa.mount
[Unit]
Description=Fex rootfs mesa layer

[Mount]
What=/usr/share/fex-emu-rootfs-layers/gentoo/images/30-mesa.sqfs
Where=/usr/share/fex-emu-rootfs-layers/gentoo/layers/30-mesa/
